# service/beringei_stats_writer/beringei_stats_writer_1year.service
[Unit]
Description= Production Beringei Stats Writer (1year)
           
[Service]
Type=simple
User=root
EnvironmentFile=/home/nms/tgnms/service/beringei_stats_writer/config/1year
ExecStart=/home/nms/tgnms/service/beringei_stats_writer/run
Restart=always
RestartSec=2s

[Install]
WantedBy=default.target
