# service/nms_mysql_writer/run
#!/bin/bash
export PORT=8086
pushd /home/nms/www
exec npm start 2>&1
